# pairdrop/rootfs/etc/services.d/pairdrop/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Pairdrop
# Take down the S6 supervision tree when Pairdrop fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Pairdrop crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Pairdrop stopped, restarting..."
